# srv/salt/gateway/modules/bacnet_stack/v1/files/supervisor_bacnet_stack.conf.tpl
[program:bacserv]
user={{ user }}
directory={{ install_dir }}
command={{ install_dir }}/bin/bacserv 123
autostart=false
